# s6/init/up
/bin/sh -c "
set -euo pipefail

# create directories
mkdir -p /etc/d2ray/logs
mkdir -p /etc/d2ray/certs
rm -rf /etc/d2ray/users

chown -R docker:docker /etc/d2ray

exec s6-envdir /var/run/s6/container_environment exec s6-setuidgid docker python3 /opt/init.py
"