# opt/crontabs/root
0 0 * * * certbot renew --post-hook "supervisorctl xray restart"