# opt/crontabs/root
* * * * * certbot renew --post-hook "supervisorctl xray restart"